# debian/preinst
#!/bin/bash

if [ install = "$1" -o upgrade = "$1" ]; then
   dpkg-divert --package sysv-freeze --add --rename \
      --divert /usr/sbin/update-rc.d.real /usr/sbin/update-rc.d
fi
